# gamepad-light-controller.service
[Unit]
Description=Gamepad Light Controller (gamepad → Zigbee lights)
After=network.target mosquitto.service zigbee2mqtt.service bluetooth.service
Requires=mosquitto.service bluetooth.service
Wants=zigbee2mqtt.service

[Service]
Type=simple
User=sparrow
Group=sparrow
WorkingDirectory=/home/sparrow/projects/zigbeendicate-sounds
ExecStart=/home/sparrow/projects/zigbeendicate-sounds/venv/bin/python3 /home/sparrow/projects/zigbeendicate-sounds/gamepad_light_controller.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Input device access
SupplementaryGroups=input bluetooth

# Environment
Environment=PYTHONPATH=/home/sparrow/projects/zigbeendicate-sounds
Environment=PYTHONDONTWRITEBYTECODE=1

[Install]
WantedBy=multi-user.target
